# debian/immich-server.service
# Please consider only editing this file through systemd drop-in files
# (e.g. `systemctl edit immich-server`). This will enable you to get updates
# to this file automatically, without your overrides being lost.

# Refer to the manpage `systemd.unit` for detailed documentation on each option.
[Unit]
Description=immich main server
Documentation=https://github.com/immich-app/immich
# The service requires network interfaces to be configured.
After=network.target
# If you want to expose the server to the network directly, you may find the following useful:
#After=network-online.target

# If immich-machine-learning, postgresql and/or redis are running on this system they should be specified as dependencies,
# such that immich-server is only started after (`After=`) those have successfully (`Requires=`) started.
#After=immich-machine-learning.service postgresql.service redis.service
#Requires=immich-machine-learning.service postgresql.service redis.service

# If running immich behind a reverse proxy - running on this system - consider ordering it before also:
#Before=theproxy.service

# Refer to the manpages `systemd.service` and `systemd.exec` for detailed documentation on each option.
[Service]
## Service execution environment
# Service configuration file.
EnvironmentFile=/etc/default/immich-server
# build info
EnvironmentFile=/usr/share/doc/immich-server/build.env
# Service execution.
Type=exec
ExecStart=/bin/sh -c ' \
  set -a; \
  . /etc/immich/server.env; \
  set +a; \
  if [ -z "$UV_THREADPOOL_SIZE" ]; then \
    CPU_CORES=$(/usr/lib/immich/server/bin/get-cpus.sh); \
    if [ "$CPU_CORES" -gt 4 ]; then \
      UV_THREADPOOL_SIZE=$CPU_CORES; \
    fi; \
  fi; \
  LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/lib/jellyfin-ffmpeg/lib \
  PATH=/usr/lib/jellyfin-ffmpeg:$PATH \
  exec node dist/main "$@"'
WorkingDirectory=/usr/lib/immich/server
# Restart on failure
Restart=on-failure
RestartSec=5s

## Permissions and directories
User=immich
Group=immich
# Read-only configuration directory (read-only mount, permissions managed by tmpfiles, warns if mode differs).
ConfigurationDirectory=immich::ro
ConfigurationDirectoryMode=0500
# Writable directories (read-write for the service, read-only for the group, reapplies permissions recursively if different).
StateDirectory=immich immich/data immich/plugins
StateDirectoryMode=0750
# File creation mask (matching modes above).
UMask=0027

## Connection and process limits
LimitNOFILE=1048576
LimitNPROC=512
#TasksMax=1024
# Customize limits to your setup and expected usage spikes.
# See: https://docs.immich.app/install/requirements/#hardware
#CPUQuota=400%
#MemoryHigh=3G
#MemoryMax=4G

## Sandboxing and isolation
# An overview is available using `systemd-analyze security immich-server`.

## Capabilities
AmbientCapabilities=
CapabilityBoundingSet=
# Uncomment these to allow immich-server to bind to ports in the range 0-1024.
#AmbientCapabilities=CAP_NET_BIND_SERVICE
#CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Process restrictions
NoNewPrivileges=true
RestrictSUIDSGID=true
SecureBits=noroot noroot-locked

# Filesystem protection
ProtectSystem=strict
ProtectHome=true
ProtectProc=noaccess
ProcSubset=pid
# Allow reading CPU info for automatic UV_THREADPOOL_SIZE calculation
BindReadOnlyPaths=/proc/cpuinfo

# Namespace isolation
PrivateTmp=disconnected
# Allow hardware-accelerated transcoding.
PrivateDevices=false
# For more security with GPU access, use DevicePolicy with DeviceAllow:
#DevicePolicy=closed
#DeviceAllow=/dev/dri rw
#DeviceAllow=/dev/nvidia* rw
PrivateIPC=true
PrivateBPF=true

# Kernel and system protection
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=strict

# Execution restrictions
# Allow for V8 JIT compiler.
MemoryDenyWriteExecute=false
SetLoginEnvironment=false

# System call and network filtering
SystemCallFilter=@system-service
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictFileSystems=@basic-api @common-block
RestrictRealtime=true

# Cleanup
RemoveIPC=true

[Install]
WantedBy=multi-user.target
